# systemd/fan.service
[Unit]
Description=PWM controlled fan.

[Service]
ExecStart/opt/FanService/bin/fanService

[Install]
WantedBy=multi-user.target